# systemd/forgejo.service
[Unit]
Description=Forgejo standalone rootless Podman compose stack
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/deploy-forgejo-standalone
ExecStart=%h/deploy-forgejo-standalone/scripts/deploy.sh
ExecStop=%h/deploy-forgejo-standalone/scripts/stop.sh
RemainAfterExit=yes
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=default.target
